# vpn-proxy.cron
31 6    * * *   root sh /opt/clash/sync_config.sh 
