# src/clash.service
[Unit]
Description=clash daemon

[Service]
Type=simple
User=root
ExecStart=/opt/clash/clash -d /opt/clash
Restart=on-failure

[Install]
WantedBy=multi-user.target
